# auto_update_github.cron
0 11 * * 6 /home/zzl/shell/auto_update_github.sh

